# upstart/armstrong-django.conf
start on starting armstrong
stop on stopping armstrong
