# upstart/armstrong_mounted.conf
start on started

script
  start wait_for_file FILE_PATH=/sites/armstrong/config/development.py
end script
