# upstart/armstrong.conf
start on file FILE_PATH=/sites/armstrong/config/development.py
stop on shutdown

pre-start script

bash << "EOF"
  mkdir -p /var/log/armstrong
  chown -R armstrong /var/log/armstrong
EOF

end script
